# crontab
0 8 * * * /usr/bin/python3 /yfay.py >> /var/log/yfay.log 2>&1
